# paper-explainer.service
[Unit]
Description=Paper Explainer Streamlit Service
After=network.target

[Service]
User=<YOUR_USER>
Group=<YOUR_USER>
WorkingDirectory=<PATH_TO_PROJECT_REPO>
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Load API keys and other secrets securely
EnvironmentFile=/etc/paper-explainer.env
ExecStart=<PATH_TO_PROJECT_REPO>/venv/bin/streamlit run streamlit_app_enhanced.py --server.port 8502 --server.address 0.0.0.0
Restart=always

[Install]
WantedBy=multi-user.target
